# systemd/goesimage.timer
[Unit]
Description=NOAA Geostationary Operational Environment Satellite Background Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
